# recipes-core/lcd-tools/lcd-tools/lcd-sync-time.service
[Unit]
Description=Sync LCD time

[Service]
Type=oneshot
ExecStart=/usr/bin/lcd-tools --sync-time

[Install]
WantedBy=basic.target
