# recipes-core/synclcdtime/synclcdtime/synclcdtime.service
[Unit]
Description=Sync LCD time

[Service]
Type=oneshot
ExecStart=/usr/bin/synclcdtime

[Install]
WantedBy=basic.target
